# script/ethercat.service.in
#
# EtherCAT master kernel modules
#

[Unit]
Description=EtherCAT Master Kernel Modules

#
# Uncomment this, if the generic Ethernet driver is used. It assures, that the
# network interfaces are configured, before the master starts.
#
#Requires=network.service # Stop master, if network is stopped
#After=network.service # Start master, after network is ready

#
# Uncomment this, if a native Ethernet driver is used. It assures, that the
# network interfaces are configured, after the Ethernet drivers have been
# replaced. Otherwise, the networking configuration tools could be confused.
#
#Before=network.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@prefix@/sbin/ethercatctl start
ExecStop=@prefix@/sbin/ethercatctl stop

[Install]
WantedBy=multi-user.target
